# kubestart.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/kube_start.sh

[Install]
WantedBy=default.target
